# sensor-serial-server.service
[Unit]
Description=Sensor Serial Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/path/to/your/project
ExecStart=/path/to/your/project/venv/bin/python3 /path/to/your/project/sensor-serial-server.py
Restart=on-failure
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target